# scripts/battery-notify/battery-notify.timer
[Unit]
Description=Run Battery Notifier every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
